# autostart/picom.desktop.in
[Desktop Entry]
Type=Application
Name=Picom (X Compositor)
GenericName=X compositor
Comment=A X compositor
TryExec=picom
Exec=picom --dbus

# turned off by default
Hidden=true
OnlyShowIn=LXQt;
X-LXQt-Module=true
